# stream.supervisor.conf.in
# The minimal Supervisord configuration for running the live stream.
[program:serve]
command=@CMAKE_INSTALL_PREFIX@/stream.sh hr
process_name=%(program_name)s
numprocs=1
directory=@CMAKE_INSTALL_PREFIX@
autostart=true
autorestart=true
user=pi

[program:nginx]
command=/usr/local/nginx/sbin/nginx -g "daemon off;"
process_name=%(program_name)s
numprocs=1
autostart=true
autorestart=true
user=root